# shinken/shinken-scheduler.conf
[program:shinken-scheduler]
command=/usr/bin/shinken-scheduler -c /etc/shinken/schedulerd.ini
user=shinken
