# shinken/shinken-poller.conf
[program:shinken-poller]
command=/usr/bin/shinken-poller -c /etc/shinken/pollerd.ini
user=shinken
