# shinken/shinken-receiver.conf
[program:shinken-receiver]
command=/usr/bin/shinken-receiver -c /etc/shinken/receiverd.ini
user=shinken
