# shinken/shinken-broker.conf
[program:shinken-broker]
command=/usr/bin/shinken-broker -c /etc/shinken/brokerd.ini
user=shinken
