# shinken/shinken-reactionner.conf
[program:shinken-reactionner]
command=/usr/bin/shinken-reactionner -c /etc/shinken/reactionnerd.ini
user=shinken
